# demo-service/supervisord.conf
[supervisord]
nodaemon=true

[program:grpc]
command=/app/demo-service -mode=grpc

[program:http]
command=/app/demo-service -mode=http